# rootfs/etc/s6-overlay/s6-rc.d/init-nodered/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Node-RED
# Configure Node-RED avant exécution
# ==============================================================================
# Assure que la configuration existe
if ! bashio::fs.directory_exists '/config/smartphoton_jkbms_bus2/'; then
    mkdir -p /config/smartphoton_jkbms_bus2/nodes \
        || bashio::exit.nok "Failed to create node-red configuration directory"
    # Copie les fichiers modèles
    cp /etc/node-red/flows.json /config/smartphoton_jkbms_bus2/
    cp /etc/node-red/settings.js /config/smartphoton_jkbms_bus2/
    # Création d'un ID de flow aléatoire
    id=$(node -e "console.log((1+Math.random()*4294967295).toString(16));")
    sed -i "s/%%ID%%/${id}/" "/config/smartphoton_jkbms_bus2/flows.json"
fi
# Vérification de version
if ! bashio::fs.file_exists '/config/smartphoton_jkbms_bus2/version'; then
    touch /config/smartphoton_jkbms_bus2/version
    VERSION=$(bashio::addon.version)
    bashio::log.info "Module Addon Version number: ${VERSION}"
    echo $VERSION >> /config/smartphoton_jkbms_bus2/version
    chmod 775 /config/smartphoton_jkbms_bus2/version
    bashio::log.info "${VERSION}"
else
    VERSIONOLD=$(cat /config/smartphoton_jkbms_bus2/version)
    VERSION=$(bashio::addon.version)
    if [ $VERSIONOLD != $VERSION ]; then
        bashio::log.info "Mise a jour de la nouvelle version"
        bashio::log.info "Module Addon Version number: ${VERSION}"
        if bashio::fs.file_exists '/config/smartphoton_jkbms_bus2/flows.backup'; then
            rm /config/smartphoton_jkbms_bus2/flows.backup
        fi
        mv /config/smartphoton_jkbms_bus2/flows.json /config/smartphoton_jkbms_bus2/flows.backup
        cp /etc/node-red/flows.json /config/smartphoton_jkbms_bus2/
    fi
fi
# Remplace le fichier flows par celui d'origine
if bashio::fs.file_exists '/config/smartphoton_jkbms_bus2/flows.backup'; then
    rm /config/smartphoton_jkbms_bus2/flows.backup
fi
mv /config/smartphoton_jkbms_bus2/flows.json /config/smartphoton_jkbms_bus2/flows.backup
cp /etc/node-red/flows.json /config/smartphoton_jkbms_bus2/

# ==============================================================================
# BUS2 PATCH: Automatiquement ajouter le préfixe B2_ aux noeuds broadcasting
# Utilise Node.js (disponible dans l'image) au lieu de python3
# ==============================================================================
bashio::log.info "Applying BUS2 prefix patch to flows.json..."
node - <<'JSEOF'
const fs = require('fs');
const f = '/config/smartphoton_jkbms_bus2/flows.json';
const BUS_PREFIX = 'B2_';
const TARGETS = ['Setup To MQTT', 'Data To MQTT', 'TRAITEMENT DES ALARMES', 'AGRÉGATION GLOBALE'];

let data;
try {
    data = JSON.parse(fs.readFileSync(f, 'utf8'));
} catch(e) {
    console.error('BUS2 PATCH ERROR: cannot read flows.json:', e.message);
    process.exit(0); // ne pas bloquer le démarrage
}

let count = 0;
for (const node of data) {
    if (node.type !== 'function') continue;
    if (!TARGETS.includes(node.name)) continue;
    let func = node.func || '';
    if (func.includes('busPrefix')) {
        console.log('Already patched: ' + node.name);
        continue;
    }

    func = '// BUS2 AUTO-PATCH\nconst busPrefix = "' + BUS_PREFIX + '";\n\n' + func;

    func = func
        .replace('var baseTopic = "BMS_" + slaveIdStr;', 'var baseTopic = busPrefix + "BMS_" + slaveIdStr;')
        .replace('const baseTopic = `BMS_${slaveIdStr}`;', 'const baseTopic = busPrefix + `BMS_${slaveIdStr}`;')
        .replace('const baseTopic = `BMS_${bmsId}`;', 'const baseTopic = busPrefix + `BMS_${bmsId}`;')
        .replace('var deviceName = "BMS_" + slaveIdStr;', 'var deviceName = busPrefix + "BMS_" + slaveIdStr;')
        .replace('const deviceName = `BMS_${slaveIdStr}`;', 'const deviceName = busPrefix + `BMS_${slaveIdStr}`;')
        .replace('const deviceName = `BMS_${bmsId}`;', 'const deviceName = busPrefix + `BMS_${bmsId}`;')
        .replace('var deviceID = "BMS_" + slaveIdStr + "_device";', 'var deviceID = busPrefix + "BMS_" + slaveIdStr + "_device";')
        .replace('const deviceID = `BMS_${bmsId}_device`;', 'const deviceID = busPrefix + `BMS_${bmsId}_device`;')
        .replace('bmsName = "BMS_master";', 'bmsName = busPrefix + "BMS_master";')
        .replace('bmsName = (bmsId === 0) ? "BMS_master" : `BMS_${bmsId}`;', 'bmsName = (bmsId === 0) ? (busPrefix + "BMS_master") : (busPrefix + `BMS_${bmsId}`);')
        .replace('const globalDeviceID = "BMS_GLOBAL_device";', 'const globalDeviceID = busPrefix + "BMS_GLOBAL_device";\n    const globalBaseTopic = busPrefix + "BMS_GLOBAL";')
        .replace('"BMS_GLOBAL/alarm/state"', '`${globalBaseTopic}/alarm/state`')
        .replace('flow.get("discoverySent_GLOBAL_ALARM")', 'flow.get(busPrefix + "discoverySent_GLOBAL_ALARM")')
        .replace('flow.set("discoverySent_GLOBAL_ALARM", true)', 'flow.set(busPrefix + "discoverySent_GLOBAL_ALARM", true)');

    node.func = func;
    count++;
    console.log('Patched: ' + node.name);
}

try {
    fs.writeFileSync(f, JSON.stringify(data));
    console.log('BUS2 patch complete: ' + count + ' nodes patched');
} catch(e) {
    console.error('BUS2 PATCH ERROR: cannot write flows.json:', e.message);
}
JSEOF

PATCHED=$(grep -c "busPrefix" /config/smartphoton_jkbms_bus2/flows.json || true)
bashio::log.info "BUS2 patch verification: ${PATCHED} occurrences of busPrefix found"
# ==============================================================================
# END BUS2 PATCH
# ==============================================================================

# Désinstallation des packages Node-RED conflictuels
cd /config || bashio::exit.nok "Could not change directory to Node-RED"
if bashio::fs.file_exists "/config/package.json"; then
    npm uninstall \
        node-red-contrib-home-assistant \
        node-red-contrib-home-assistant-llat \
        node-red-contrib-home-assistant-ws \
        || bashio::log.warning "Failed un-installing conflicting packages"
fi
